# inst/tgester.service.template
[Unit]
Description=Telegram News Summary
After=network.target

[Service]
Type=oneshot
User={{USER}}
WorkingDirectory={{DEPLOY_DIR}}
ExecStart={{DEPLOY_DIR}}/.venv/bin/tgester summarise \
    --config {{DEPLOY_DIR}}/config.yaml \
    --env {{DEPLOY_DIR}}/.env
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
